# deploy/kat-portal-backup.service
# Oneshot backup, dipicu oleh kat-portal-backup.timer.
# Salin ke /etc/systemd/system/ bersama .timer-nya.

[Unit]
Description=Backup DB Portal Panitia KAT ITB 2026

[Service]
Type=oneshot
User=kat-portal
Group=kat-portal
EnvironmentFile=/etc/kat-portal/api.env
ExecStart=/opt/kat-portal-api/server/deploy/backup.sh
